# provisioning/roles/common/files/daemontools.conf
# This task is run on startup to start daemontools

description "start svscanboot on startup to run daemontools"

start on startup

task
exec /usr/bin/svscanboot &
